# deploy/auditmaton-ads-huey-staging.service
# deploy/auditmaton-ads-huey-staging.service
# systemd unit file for the Auditmaton: Ads Huey background worker on STAGING.
# Processes queued chart generation and PDF exports in a separate process from the
# Gunicorn web server. Distinct from the production unit (deploy/auditmaton-ads-huey.service)
# because staging lives under ~/apps/staging/auditmaton-ads while prod will
# live under ~/apps/auditmaton-ads.
#
# Install:
#   sudo cp deploy/auditmaton-ads-huey-staging.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable auditmaton-ads-huey-staging
#   sudo systemctl start auditmaton-ads-huey-staging
#
# Logs:
#   sudo journalctl -u auditmaton-ads-huey-staging -f

[Unit]
Description=Auditmaton: Ads Staging Huey Worker
After=network.target postgresql.service

[Service]
User=anniecushing
Group=www-data
WorkingDirectory=/home/anniecushing/apps/staging/auditmaton-ads
Environment="PATH=/home/anniecushing/apps/staging/auditmaton-ads/venv/bin"
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/anniecushing/apps/staging/auditmaton-ads/venv/bin/huey_consumer \
    services.task_queue.huey \
    -w 4 \
    -k process \
    --logfile=/home/anniecushing/apps/staging/auditmaton-ads/logs/huey.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
